# .config/systemd/user/nt-pipewire.service
[Unit]
Description=Noisetorch using pipewire

[Service]
ExecStart=/home/lfron/nt-pw-poc

[Install]
WantedBy=sound.target
